# deploy/forexbot.service
[Unit]
Description=Forex Trading Bot (MetaTrader 5 under Wine)
After=network-online.target xvfb99.service
Wants=network-online.target
Requires=xvfb99.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/forexbot
ExecStart=/home/ubuntu/forexbot/deploy/run_bot.sh
Restart=always
RestartSec=20
TimeoutStopSec=30
KillMode=mixed

[Install]
WantedBy=multi-user.target
